# system_files/etc/systemd/system/var-mnt-tuxbox-music.automount
[Unit]
Description=Automount Music NFS Share
ConditionPathExists=/var/mnt/tuxbox

[Automount]
Where=/var/mnt/tuxbox/music

[Install]
WantedBy=multi-user.target
