# rootfs/usr/lib/systemd/system/var-mnt-tuxbox-learning.automount
[Unit]
Description=Automount Learning NFS Share
ConditionPathExists=/var/mnt/tuxbox

[Automount]
Where=/var/mnt/tuxbox/learning

[Install]
WantedBy=multi-user.target
