# rootfs/usr/lib/systemd/system/var-mnt-tuxbox-photos.automount
[Unit]
Description=Automount Photos NFS Share
ConditionPathExists=/var/mnt/tuxbox

[Automount]
Where=/var/mnt/tuxbox/photos

[Install]
WantedBy=multi-user.target
