# system_files/usr/lib/systemd/system/var-mnt-tuxbox-television.automount
[Unit]
Description=Automount Television NFS Share
ConditionPathExists=/var/mnt/tuxbox

[Automount]
Where=/var/mnt/tuxbox/television

[Install]
WantedBy=multi-user.target
