# system_files/usr/lib/systemd/system/var-mnt-tuxbox-movies.automount
[Unit]
Description=Automount Movies NFS Share
ConditionPathExists=/var/mnt/tuxbox

[Automount]
Where=/var/mnt/tuxbox/movies

[Install]
WantedBy=multi-user.target
